# meta/recipes-core/base-files/files/genio-1200-evk/usbmass.service
[Unit]
Description=USB Mass Storage Gadget

[Service]
ExecStart=/lib/systemd/usbmass.sh --start
RemainAfterExit=true
ExecStop=/lib/systemd/usbmass.sh --stop

[Install]
WantedBy=multi-user.target
